# udev/hdmi.rules
KERNEL=="card0", ACTION=="change", RUN+="/home/ja/bin/hdmi.sh"
